# service/privoxy/run
#!/bin/sh
exec su-exec privoxy privoxy --no-daemon
